# ansible/roles/nas/templates/backup.timer.j2
[Unit]
Description=Runs autorestic backup daily

[Timer]
OnCalendar=*-*-* {{ nas_backup_time }}
Persistent=true

[Install]
WantedBy=timers.target
