# backups/autorestic/backup.service
[Unit]
Description=Runs autorestic backup
ConditionPathExists=/mnt/kd

[Service]
ExecStart=/home/debian/scripts/backup-script
