# backups/autorestic/backup.timer
[Unit]
Description=Runs autorestic backup daily

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true

[Install]
WantedBy=timers.target
