# ansible/roles/arch_repository/templates/update-aur.service.j2
[Unit]
Description=Updates 'aur' repo
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/aura sync
CPUWeight=1
IOWeight=1
